# deploy/orangepi/palmgate-device.service
[Unit]
Description=PalmGate Device Runtime
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/palmgate
Environment=DEVICE_RUNTIME_ENABLED=1
Environment=CAMERA_SOURCE=usb
Environment=CAMERA_DEVICE_INDEX=0
ExecStart=/opt/palmgate/.venv/bin/python -m app.device_runtime
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
