# deploy/orangepi/palmgate-api.service
[Unit]
Description=PalmGate API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/palmgate
ExecStart=/opt/palmgate/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
